# src/gpiod.service
[Unit]
Description=Simplified GPIO listener
After=syslog.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi
ExecStart=/home/pi/gpio-listener.py
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target


